# exynos5-hwcomposer.service
[Unit]
Description=Exynos5 hwcomposer
BindsTo=dev-fb0.device

[Service]
ExecStart=/usr/sbin/exynos5-hwcomposer
